# agg/MODULE/files/config/iecar2vr.crontab
####################################################
# Run the available Aggregators
####################################################
* * * * * /usr/revo/datafeed/aggregator_config/bin/spawn_aggregator 1>>/usr/revo/datafeed/log/iecar2vr.spawn_aggregator.log 2>&1

* * * * * /usr/revo/datafeed/aggregator_config/bin/decode_wrapper P03 ODD  PRE  1>>/usr/revo/datafeed/log/decode_wrapper.P03.ODD.PRE.log   2>&1
* * * * * /usr/revo/datafeed/aggregator_config/bin/decode_wrapper P03 ODD  POST 1>>/usr/revo/datafeed/log/decode_wrapper.P03.ODD.POST.log  2>&1
* * * * * /usr/revo/datafeed/aggregator_config/bin/decode_wrapper P03 EVEN PRE  1>>/usr/revo/datafeed/log/decode_wrapper.P03.EVEN.PRE.log  2>&1
* * * * * /usr/revo/datafeed/aggregator_config/bin/decode_wrapper P03 EVEN POST 1>>/usr/revo/datafeed/log/decode_wrapper.P03.EVEN.POST.log 2>&1

